# floreria-bot.service
# /etc/systemd/system/floreria-bot.service
# Jardín RoCe — WhatsApp Bot optimizado para GCP e2-micro
# Instalación:
#   sudo cp floreria-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable floreria-bot
#   sudo systemctl start floreria-bot

[Unit]
Description=Jardín RoCe WhatsApp Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cantedavid00
WorkingDirectory=/home/cantedavid00/jardin-roce-bot

# Node.js: heap limitado a 380MB, GC agresivo
# NOTA: .env.local lo carga el bot via dotenv, no via EnvironmentFile
Environment="NODE_OPTIONS=--max-old-space-size=380 --max-semi-space-size=32"

# PATH completo para que systemd encuentre npm
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

ExecStart=npm run bot:prod

# ═══ Reinicio automático ═══
Restart=on-failure
RestartSec=10s

# ═══ Timeout al hacer stop (evita el "stop-sigterm timed out") ═══
TimeoutStopSec=15s
KillSignal=SIGTERM
SendSIGKILL=yes

# ═══ Límites de memoria (cgroups) — OOM killer controlado ═══
# 650M deja ~350M para el OS en 1GB total
MemoryMax=650M
MemoryHigh=550M

# ═══ Límites de CPU ═══
CPUQuota=80%

# ═══ Logs ═══
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
