# floreria-bot.service
# /etc/systemd/system/floreria-bot.service
# Jardín RoCe — WhatsApp Bot optimizado para GCP e2-micro
# Instalación:
#   sudo cp floreria-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable floreria-bot
#   sudo systemctl start floreria-bot

[Unit]
Description=Jardín RoCe WhatsApp Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=antedavid00
WorkingDirectory=/home/antedavid00/jardin-roce-bot

# Cargar variables de entorno
EnvironmentFile=/home/antedavid00/jardin-roce-bot/.env.local

# Node.js: heap limitado a 380MB, GC agresivo
Environment=NODE_OPTIONS=--max-old-space-size=380 --optimize-for-size --max-semi-space-size=32

# Comando: usa el script bot:prod para heredar NODE_OPTIONS
ExecStart=npm run bot:prod

# ═══ Reinicio automático ═══
Restart=on-failure
RestartSec=10s

# ═══ Timeout al hacer stop (evita el "stop-sigterm timed out") ═══
TimeoutStopSec=15s
KillSignal=SIGTERM
SendSIGKILL=yes

# ═══ Límites de memoria (cgroups) — OOM killer controlado ═══
# 650M deja ~350M para el OS en 1GB total
MemoryMax=650M
MemoryHigh=550M

# ═══ Límites de CPU ═══
CPUQuota=80%

# ═══ Logs ═══
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
